# supervisord/supervisord.conf
[unix_http_server]
file=/opt/logs/supervisord/supervisor.sock 
chmod=0700 

[supervisord]
logfile=/opt/logs/supervisord/supervisord.log 
pidfile=/opt/logs/supervisord/supervisord.pid
nodaemon=false               ; start in foreground if true; default false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///opt/logs/supervisord/supervisor.sock

[program:zoro]
user=root
command=/opt/work/zoro/src/run.sh
autostart=true
autorestart=true
stdout_logfile=/opt/logs/supervisord/zoro.log
stderr_logfile=/opt/logs/supervisord/zoro.log

[program:nginx]
user=root
command=/usr/local/nginx/nginx
numprocs=2
autostart=true
autorestart=true

